# examples/FileShare.service
[Unit]
Description=FileShare Website
After=network.target

[Service]
Type=simple
ExecStart=/path/to/start.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
